# systemd/inference-worker.service
[Unit]
Description=Inference Worker
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/slm-deploy/quickstart/workers/inference-worker
Environment=III_URL=ws://10.0.0.174:49134
ExecStart=/home/ubuntu/slm-deploy/quickstart/workers/inference-worker/venv/bin/python inference_worker.py
Restart=always

[Install]
WantedBy=multi-user.target